# Service/webServer.service
[Unit]
Description="Web Service"

[Service]
ExecStart=<executable>

[Install]
WantedBy=multi-user.target