# Service/webServer.service
[Unit]
Description="Web Service"

[Service]
ExecStart={{launchCommand}}

[Install]
WantedBy=multi-user.target